# updater/crontab
0 0 * * * /updater/update.sh > /proc/1/fd/1 2>&1
